# config/playit.service
# LEGACY — This file is no longer installed by bin/setup-playit.sh.
#
# When installing via the official Playit APT package, the vendor ships its
# own unit at /usr/lib/systemd/system/playit.service. Copying a custom file
# into /etc/systemd/system/playit.service overwrites the vendor unit and can
# cause agent identity mismatches (duplicate secret keys, tunnels online but
# receiving no traffic).
#
# setup-playit.sh now writes a drop-in override instead:
#   /etc/systemd/system/playit.service.d/override.conf
#
# This file is kept only as a historical reference and is NOT used by any
# script in this repository.
#
# ─────────────────────────────────────────────────────────────────────────────
[Unit]
Description=Playit.gg Tunnel for Minecraft
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=playit
Group=playit
WorkingDirectory=/opt/playit
ExecStart=/opt/playit/playit
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
